# systemd/user/runonce@.service
[Unit]
Description=Runs %i once
ConditionPathExists=%h/.runonce/queue/%i

[Service]
Type=oneshot
ExecStart=/usr/sbin/runonce %i
